# supervisor/chatwoot.conf
[program:chatwoot]
priority=10
directory=/app
command=bundle exec rails s -p 3000 -b 0.0.0.0
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
