# supervisor/sidekiq.conf
[program:sidekiq]
priority=10
directory=/app
command=bundle exec sidekiq -C config/sidekiq.yml
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
